# gunicorn/truckman_gunicorn.socket
[Unit]
Description=gunicorn socket
[Socket]
ListenStream=/run/truckman_gunicorn.sock
[Install]
WantedBy=sockets.target